# lucas-live-trading/deploy/systemd/tradingbot-bot.service.template
[Unit]
Description=TradingBot live trading engine (bot.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__REPO_DIR__
ExecStart=__REPO_DIR__/.venv/bin/python lucas-live-trading/bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Hardening: the bot only needs to read/write within the repo dir.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
